# deploy/systemd/banker-radar-daily-review.service
[Unit]
Description=Banker Radar daily review Telegram push
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/banker-radar
EnvironmentFile=-/root/banker-radar/.env
ExecStart=/usr/bin/flock -n /run/banker-radar-daily-review.lock /usr/local/bin/banker-radar telegram-review --period yesterday --config /root/banker-radar/configs/radar.yaml
