# deploy/systemd/banker-radar-tracking.timer
[Unit]
Description=Run Banker Radar signal tracking every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
